# sonos-flac.timer
[Unit]
Description=Run Sonos FLAC Conversion nightly

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
